# deploy/ansible/roles/app/templates/app.service.j2
[Unit]
Description={{ app_name }} Spring Boot Application
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory=/opt/{{ app_name }}

ExecStart=/usr/bin/java \
    -Xms{{ java_heap_min }} \
    -Xmx{{ java_heap_max }} \
    -XX:+UseZGC \
    -XX:+DisableAttachMechanism \
    -XX:+HeapDumpOnOutOfMemoryError \
    -XX:HeapDumpPath=/var/log/{{ app_name }}/heapdump.hprof \
    -Djava.security.egd=file:/dev/urandom \
    -jar /opt/{{ app_name }}/{{ app_name }}.jar \
    --spring.config.location=/opt/{{ app_name }}/application.properties

# Security: Prevent core dumps that may contain sensitive data
LimitCORE=0

# Security: Restrict filesystem access
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/{{ app_name }} /var/log/{{ app_name }}
PrivateTmp=true

# Security: Restrict capabilities
NoNewPrivileges=true
CapabilityBoundingSet=

SuccessExitStatus=143
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ app_name }}

[Install]
WantedBy=multi-user.target
